# installer_rlabs/packages/uwsgi/uwsgi.service
[Unit]
Description=uWSGI Emperor service
After=syslog.target


[Service]
ExecStart=/usr/local/bin/uwsgi --emperor /etc/uwsgi/sites/web2py.ini
Restart=always
KillSignal=SIGQUIT
Type=notify
NotifyAccess=all


[Install]
WantedBy=multi-user.target


